# ghost_upstart.conf
# ghost

# author "Karn Kaul"
# description "Ensures Ghost is always running"

start on startup

script
	cd /root/ghost
	npm start
end script